# crontab
* * * * * /bin/bash ~/deployment-scripts/utils/repo-check.sh > ~/deployment-scripts.log