# crontab
* * * * * /bin/bash ~/deployment-scripts/utils/update-repository.sh > ~/deployment-scripts.log
