# crontab
* * * * * /bin/bash ~/deployment-scripts/utils/update-repository.sh > /var/log/deployment-scripts.log
