# crontab
* * * * * /bin/bash ~/deployment-scripts/utils/check-for-changes.sh > ~/deployment-scripts.log
